# mkdocs/rootfs/etc/services.d/mkdocs/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Finish the MkDocs service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "MkDocs service stopped"