# mkdocs/rootfs/etc/services.d/rebuild-api/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting MkDocs rebuild API on port 8081"
exec python3 -u /usr/local/bin/rebuild_api_server.py --port 8081
