# mkdocs/rootfs/etc/services.d/mkdocs-api/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# MkDocs API service finish script
# ==============================================================================

bashio::log.info "MkDocs API service stopped"