# mkdocs/rootfs/etc/services.d/mkdocs-api/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# MkDocs API service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare enable_api

# Get configuration
enable_api=$(bashio::config 'enable_api' 'true')

if [[ "${enable_api}" == "true" ]]; then
    bashio::log.info "Starting MkDocs API service"
    # Start the API server
    exec /usr/bin/mkdocs-api-server
else
    bashio::log.info "MkDocs API service disabled"
    # Sleep indefinitely to keep the service running but inactive
    exec sleep infinity
fi